# system/supervisor/mongodb.ini
[program:mongodb]
autorestart=true
command=/usr/bin/mongod --config /etc/mongodb.conf
priority=995
redirect_stderr=true
user=mongodb
