# system/supervisor/redis.ini
[program:redis_castle-db0]
autorestart=true
command=/opt/redis-2.6.12/redis-server /etc/redis/castle-db0.conf
priority=996
redirect_stderr=true
user=redis
