# system/supervisor/stud.ini
[program:stud]
command=/usr/local/bin/stud --config=/etc/stud/stud.cfg
user=www
autostart=true
priority=998
